# systemd/machinevisionflow-backend.service
[Unit]
Description=Machine Vision Flow - Python Backend
Documentation=https://github.com/your-repo/MachineVisionFlow
After=network.target

[Service]
Type=forking
User=cnc
Group=cnc
WorkingDirectory=/home/cnc/machine-vision/backend/python-backend

# Runtime directories (systemd creates these with correct permissions)
RuntimeDirectory=machinevisionflow
LogsDirectory=machinevisionflow

# Environment variables for production paths
Environment="LOG_DIR=/var/log/machinevisionflow"
Environment="RUN_DIR=/run/machinevisionflow"
Environment="MV_CONFIG_FILE=/home/cnc/machine-vision/backend/python-backend/config.yaml"

# Start command
ExecStart=/home/cnc/machine-vision/backend/python-backend/venv/bin/python main.py
ExecStop=/bin/kill -TERM $MAINPID

# PID file
PIDFile=/run/machinevisionflow/backend.pid

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Restart policy
Restart=on-failure
RestartSec=10s

# Logging
StandardOutput=append:/var/log/machinevisionflow/backend.log
StandardError=append:/var/log/machinevisionflow/backend.log

[Install]
WantedBy=multi-user.target
